# templates/automount_file.automount.j2
[Unit]
Description={{ item.description }}

[Automount]
Where={{ item.mount }}
TimeoutIdleSec={{ item.timeout_idle }}

[Install]
WantedBy={{ item.wantedby }}
